# packages/thorch-kde-defaults/payload/usr/lib/systemd/user/thorch-audio-init.service
[Unit]
Description=Initialize Thorch speaker audio routing
After=pipewire.service wireplumber.service pipewire-pulse.service

[Service]
Type=oneshot
ExecStart=/usr/bin/thorch-audio-init

[Install]
WantedBy=default.target

