# packages/thorch-bsp/payload/usr/lib/systemd/system/thorch-debug-report.service
[Unit]
Description=Write Thorch boot diagnostics to /boot
After=systemd-udevd.service NetworkManager.service display-manager.service
Wants=NetworkManager.service
ConditionPathIsMountPoint=/boot

[Service]
Type=oneshot
ExecStart=/usr/bin/thorch-debug-report

[Install]
WantedBy=multi-user.target
