# packages/thorch-kde-defaults/payload/usr/lib/systemd/system/thorch-session-recovery.service
[Unit]
Description=Recover from failed Thorch session switches
After=display-manager.service systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/bin/thorch-session-recovery

[Install]
WantedBy=multi-user.target
